# deploy/saturn-events.service
# saturn-events — Discord event deliverer (deliverer/deliverer.mjs).
# Installed to /etc/systemd/system by deploy.sh. Polls the app's
# /api/event-subscriptions feed, holds one Discord Gateway websocket per
# distinct bot token, POSTs matched mentions to /api/events.
[Unit]
Description=Saturn Discord event deliverer
Wants=network-online.target
After=network-online.target saturn.service

[Service]
User=saturn
Group=saturn
EnvironmentFile=/etc/saturn/saturn.env
Environment=SATURN_BASE_URL=http://127.0.0.1:3000
WorkingDirectory=/srv/saturn/app
ExecStart=/usr/bin/env node /srv/saturn/app/deliverer/deliverer.mjs
Restart=always
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=multi-user.target
