# deploy/saturn.service
[Unit]
Description=Saturn — Next.js app (saturn.wilcus.com)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
User=saturn
Group=saturn
WorkingDirectory=/srv/saturn/app
EnvironmentFile=/etc/saturn/saturn.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node node_modules/.bin/next start -p 3000 -H 127.0.0.1
Restart=always
RestartSec=3s
# a wedged shutdown must not stall deploys
TimeoutStopSec=15

# hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
# the only path `next start` writes (image/fetch cache)
ReadWritePaths=/srv/saturn/app/.next/cache
# reach the rootless Podman socket: group membership for the 0750 socket dir...
SupplementaryGroups=sandboxes
# ...and connect() writes the socket inode, which ProtectSystem=strict blocks
# without an explicit write grant on the containing path.
ReadWritePaths=/run/sandboxes

MemoryMax=2G
MemoryHigh=1536M
TasksMax=256

[Install]
WantedBy=multi-user.target
